# default/.config/systemd/user/newsboat.service
[Unit]
Description=Newsboat feed reload service

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload

[Install]
WantedBy=default.target
