# home/.config/systemd/user/huewarm.service
[Unit]
Description=Hue light color temperature service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mullvad_hue
