# home/.config/systemd/user/libreddit.service
[Unit]
Description=libreddit daemon
After=network.service

[Service]
ExecStart=/usr/bin/libreddit

[Install]
WantedBy=default.target
