# home/.config/systemd/user/mailsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/frank/.local/bin/mailsync
