# home/.config/systemd/user/mailsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=30
OnUnitActiveSec=10m
Unit=mailsync.service

[Install]
WantedBy=timers.target
