# home/.config/systemd/user/newsboat.timer
[Unit]
Description=Newsboat feed reload timer

[Timer]
OnBootSec=30
OnUnitActiveSec=5m
Unit=newsboat.service

[Install]
WantedBy=timers.target
