# home/.config/systemd/user/redlib.service
[Unit]
Description=redlib daemon
After=network.service

[Service]
ExecStart=/usr/bin/redlib

[Install]
WantedBy=default.target
