# home/.config/systemd/user/trash-empty.service
[Unit]
Description=Clear trash of old files

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty 32
